# client/crontab
0 * * * * /bin/sh 
/Users/ethanguo/chargepoint/sync/chargepoint_sync.command
